# root/etc/cron.daily/check_pass.cron
#! /bin/sh

/etc/e-smith/events/actions/user-passwd-expiration
